# service/svscan/run
#!/bin/sh

/usr/sbin/eecloud -c /etc/eecloud/eecloud.conf
